# scripts/agent-comm-port8000.service
[Unit]
Description=Agent Communication Server (Port 8000)
Documentation=https://github.com/moai/agent-comm
After=network-online.target postgresql.service
Wants=network-online.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/agent_com

# Environment
Environment="PATH=/home/ubuntu/.local/bin:/usr/bin:/bin"
Environment="CONFIG_PATH=/home/ubuntu/agent_com/config.production.json"
Environment="PYTHONUNBUFFERED=1"

# Process management
ExecStart=/home/ubuntu/.local/bin/uv run python -m communication_server.main
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/agent_com
ReadWritePaths=/var/log/agent-comm
LogLevel=info

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

# Logging
StandardOutput=append:/var/log/agent-comm/output.log
StandardError=append:/var/log/agent-comm/error.log
SyslogIdentifier=agent-comm

[Install]
WantedBy=multi-user.target
